# scripts/upstart/node-app.conf
description "node app"

start on vagrant-ready or runlevel [2345]
stop on runlevel [016]

respawn

exec /usr/bin/node /application/bin/www >> /var/log/node-app.log 2>&1

